# tests/benchmark/scores/TEST-001_20260502_044744_score.json
{
  "case_id": "TEST-001",
  "threat_type": "ransomware",
  "scores": {
    "ioc_precision": 0.0,
    "ioc_recall": 0.0,
    "ioc_f1": 0.0,
    "section_completeness": 0.0,
    "verdict_accuracy": 0.0,
    "overall": 0.0
  },
  "ioc_details": [
    {
      "type": "process",
      "value": "invoice.exe",
      "confidence": "high",
      "found": false,
      "found_in": ""
    },
    {
      "type": "process",
      "value": "powershell.exe",
      "confidence": "high",
      "found": false,
      "found_in": ""
    },
    {
      "type": "process",
      "value": "cmd.exe",
      "confidence": "medium",
      "found": false,
      "found_in": ""
    },
    {
      "type": "file",
      "value": "invoice.exe",
      "confidence": "high",
      "found": false,
      "found_in": ""
    },
    {
      "type": "registry_key",
      "value": "HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run",
      "confidence": "high",
      "found": false,
      "found_in": ""
    },
    {
      "type": "technique",
      "value": "T1547.001",
      "confidence": "high",
      "found": false,
      "found_in": ""
    },
    {
      "type": "technique",
      "value": "T1059.001",
      "confidence": "high",
      "found": false,
      "found_in": ""
    },
    {
      "type": "technique",
      "value": "T1486",
      "confidence": "high",
      "found": false,
      "found_in": ""
    }
  ],
  "sections_found": [],
  "sections_missing": [
    "Executive Summary",
    "Timeline of Events",
    "Indicators of Compromise",
    "Persistence Mechanisms",
    "Recommendations",
    "Evidence References"
  ],
  "verdict_keywords_found": [],
  "verdict_keywords_missing": [
    "ransomware",
    "invoice",
    "powershell",
    "persistence",
    "run key",
    "encrypt"
  ],
  "flags": {
    "persistence_detected": false,
    "lateral_movement_detected": false,
    "timestomp_detected": false,
    "code_injection_detected": false
  }
}
